# agent-fleet/daemon/packaging/linux/agent-fleet-daemon.service
[Unit]
Description=Agent-Fleet Daemon
Documentation=https://github.com/bringie/vault-rag
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentfleet
Group=agentfleet
WorkingDirectory=/opt/agent-fleet
EnvironmentFile=/etc/agent-fleet/daemon.env
ExecStart=/usr/bin/env node /opt/agent-fleet/bin/daemon.js \
    --hub  ${AGENT_FLEET_HUB} \
    --token ${AGENT_FLEET_TOKEN} \
    --host-name ${AGENT_FLEET_HOST_NAME} \
    --state-dir ${AGENT_FLEET_STATE_DIR}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-fleet-daemon

# Hardening — daemon only needs to read its source + write state-dir + run claude.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/agent-fleet
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
